# daemon/systemd/scaffold-daemon.service
[Unit]
Description=Scaffold Go Daemon
After=network-online.target scaffold-signal-cli.service
Wants=network-online.target
Requires=scaffold-signal-cli.service

[Service]
Type=simple
WorkingDirectory=%h/Builds/scaffold/daemon
EnvironmentFile=%h/Builds/scaffold/daemon/.env
Environment=PATH=/home/mikekey/.bun/bin:/home/mikekey/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/Builds/scaffold/daemon/bin/scaffold-daemon
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
